# roles/keystone/templates/etc/cron.d/drop-expired-keystone-tokens
PATH=/usr/local/bin/

# Drop the tokens every 5 minutes
*/5 * * * * root keystone-manage token_flush
